# deploy/agentsd.service
[Unit]
Description=local-agents control plane (agentsd)
After=network-online.target xvfb.service
Wants=network-online.target
# Start virtual display when available (headed Playwright)
Wants=xvfb.service

[Service]
Type=simple
# Prefer a dedicated user in production:
# User=agents
# Group=agents
EnvironmentFile=-/etc/agentsd/env
Environment=HOME=/root
Environment=DISPLAY=:99
Environment=USE_GKE_GCLOUD_AUTH_PLUGIN=True
Environment=CLOUDSDK_CORE_DISABLE_PROMPTS=1
Environment=PATH=/root/.local/bin:/root/.grok/bin:/root/.opencode/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/var/lib/local-agents
ExecStart=/usr/local/bin/agentsd serve --config /etc/agentsd/config.toml
Restart=on-failure
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
